# jupyterhub2/podman.socket
[Unit]
Description=Podman API Socket

[Socket]
ListenStream=%h/podman.sock
SocketMode=0660

[Install]
WantedBy=sockets.target
