# jupyterhub2/podman.socket
[Unit]
Description=Podman API Socket

[Socket]
ListenStream=/workdir/home/podman.sock
SocketMode=0660

[Install]
WantedBy=sockets.target
